# Hassamole/rootfs/etc/services.d/postgres/run
#!/usr/bin/with-contenv bashio
#shellcheck shell=bash

bashio::log.info "Starting postgres..."
s6-setuidgid postgres postgres
